# media-stack/systemd/glacier-watch.service
[Unit]
Description=Surveillance Glacier (upload/restauration) -> ntfy
After=jellyfin.service

[Service]
Type=oneshot
ExecStart=/srv/media-stack/glacier-watch.sh
